# deploy/recorder/birdcam-recorder.service
[Unit]
Description=birdcam recorder (RTSP ingest + clip recording + viewer backend + stream)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yourUser
WorkingDirectory=/path/to/birdcam
EnvironmentFile=/etc/birdcam/env
ExecStart=/path/to/birdcam/.venv/bin/birdcam-recorder
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
